# config/systemd/zds-clear-sessions@.timer
[Unit]
Description=Clear Django Sessions Timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
